# credential-proxy/units/credential-proxy-iptables.service
[Unit]
Description=Restore iptables rules for credential-proxy
DefaultDependencies=no
Before=network-pre.target credential-proxy.service
Wants=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c '/usr/sbin/iptables-restore --noflush /etc/credential-proxy/iptables.rules'
ExecStart=/bin/sh -c 'if [ -f /etc/credential-proxy/iptables.v6.rules ]; then /usr/sbin/ip6tables-restore --noflush /etc/credential-proxy/iptables.v6.rules; fi'
ExecStop=/bin/sh -c 'if [ -x /etc/credential-proxy/iptables.remove.sh ]; then /etc/credential-proxy/iptables.remove.sh; fi'
ExecStop=/bin/sh -c 'if [ -x /etc/credential-proxy/iptables.v6.remove.sh ]; then /etc/credential-proxy/iptables.v6.remove.sh; fi'

[Install]
WantedBy=multi-user.target
