# credential-proxy/units/credential-proxy.service
[Unit]
Description=Credential Proxy (multi-service, SNI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw-secrets
Group=openclaw-secrets
ExecStart=/usr/bin/python3 /opt/credential-proxy/credential_proxy.py --config /etc/openclaw-secrets/credential-proxy.yaml
Restart=on-failure
RestartSec=1s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/credential-proxy /etc/openclaw-secrets/certs/hosts
ReadOnlyPaths=/etc/openclaw-secrets
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
